# swoole.ini
[program:swoole]
command=php /Users/heping/PhpstormProjects/workspace/fengma-fitness-server/bin/fengma-fitness start
user=heping
redirect_stderr=true
stdout_logfile=/Users/heping/PhpstormProjects/workspace/fengma-fitness-server/app/Logs/swoolele.log
autostart=true
auturestart=true
numprocs=1